# packages/gateway/gateway.socket
[Unit]
Description=Socket for the HTTPS gateway

[Socket]
ListenStream=443
TCPCongestion=bbr

[Install]
WantedBy=sockets.target

# vim: ft=systemd ts=4 sw=4 sts=4 et

